# deploy/family-bot.service
[Unit]
Description=Family Bot - Silat Al-Rahm Reminder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=familybot
WorkingDirectory=/opt/family-bot
ExecStart=/opt/family-bot/.venv/bin/python -m bot.main
Restart=on-failure
RestartSec=10
EnvironmentFile=/opt/family-bot/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
